# field-trainer.service
[Unit]
Description=Field Trainer Circuit Training System
Documentation=file:///opt/field-trainer/README.md
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/field-trainer/app
Environment=PYTHONPATH=/opt/field-trainer/app
Environment=FIELD_TRAINER_CONFIG=/opt/field-trainer/config
Environment=FIELD_TRAINER_DATA=/opt/field-trainer/data
ExecStart=/usr/bin/python3 /opt/field-trainer/app/field_trainer_main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
KillMode=mixed

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/field-trainer/data /opt/field-trainer/config

# Resource limits
LimitNOFILE=1024
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=field-trainer

[Install]
WantedBy=multi-user.target